# templates/com.soundcore-ai-recorder.daily.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.soundcore-ai-recorder.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/absolute/path/to/soundcore-ai-recorder/scripts/run_daily.py</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/absolute/path/to/soundcore-ai-recorder/logs/daily.out.log</string>
  <key>StandardErrorPath</key>
  <string>/absolute/path/to/soundcore-ai-recorder/logs/daily.err.log</string>
</dict>
</plist>

